# systemd/soci-video-cdn.service
[Unit]
Description=Application soci-avatar-cdn

[Service]
Environment=APP_PORT=4202

[Install]
WantedBy=multi-user.target
